# host-work/config/autostart/papis-tmsu.desktop
[Desktop Entry]
Name=papis-tmsu-vfs
Exec=/usr/local/bin/papis-tmsu-vfs.sh
Type=Application
