# tag-qubes/config/autostart/ssh-agent.desktop
[Desktop Entry]
Name=ssh-agent
Exec=/rw/config/rc.local; ssh-agent
Type=Application
