# config/autostart/ssh-add.desktop
[Desktop Entry]
Name=ssh-agent
Exec=ssh-agent
Type=Application
